# supervisord.conf
[supervisord]

nodaemon=true

[program:gitbook]

command=gitbook serve /gitbook
